# root/etc/services.d/lidarr/run
#!/usr/bin/with-contenv bash

umask 022

cd /opt/lidarr || exit

exec \
	s6-setuidgid abc mono --debug Lidarr.exe \
	-nobrowser -data=/config

